# build/s6-services/etc/s6-overlay/s6-rc.d/01-set-perms/up
#!/command/execlineb -P 
bash /etc/s6-overlay/s6-rc.d/01-set-perms/set-permissions.sh
