# build/s6-services/etc/s6-overlay/s6-rc.d/10-start-console/up
#!/command/execlineb -P
s6-setuidgid svnedge
cd /home/svnedge/csvn
bin/csvn start
